# etc/init.d/pulseaudio-dirs
#!/bin/sh /etc/rc.common
# Copyright (C) 2011 OpenWrt.org

START=98

start() {
	[ -d /var/lib/pulse ] || {
		mkdir -m 0755 -p /var/lib
        ln -s /mnt/storage/pulse/ /var/lib/pulse
	}
}
